# target/supervisor_restic.ini
[supervisord]
nodaemon=true

[program:cron]
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
command=/go/bin/start_cron.sh

[program:restic_backup]
startsecs=0
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/go/bin/restic backup /data %(ENV_RESTIC_BACKUP_OPTIONS)s

[program:restic_cleanup]
startsecs=0
autostart=false
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/go/bin/restic forget --keep-daily %(ENV_RESTIC_CLEANUP_KEEP_DAILY)s --keep-weekly %(ENV_RESTIC_CLEANUP_KEEP_WEEKLY)s --keep-monthly %(ENV_RESTIC_CLEANUP_KEEP_MONTHLY)s --keep-yearly %(ENV_RESTIC_CLEANUP_KEEP_YEARLY)s --prune %(ENV_RESTIC_CLEANUP_OPTIONS)s

